# council/lord-chamberlain/lord-chamberlain.timer
[Unit]
Description=Lord Chamberlain — triage poll every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=lord-chamberlain.service

[Install]
WantedBy=timers.target
